# supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

# ──────────────────────────────────────────────
# PostgreSQL
# ──────────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=1
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ──────────────────────────────────────────────
# MongoDB
# ──────────────────────────────────────────────
[program:mongodb]
command=mongod --dbpath /var/lib/mongodb --bind_ip_all
autostart=true
autorestart=true
priority=1
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ──────────────────────────────────────────────
# FastAPI — attend que les DBs soient prêtes
# ──────────────────────────────────────────────
[program:fastapi]
command=/bin/bash -c "
  echo '[fastapi] Waiting for PostgreSQL...';
  until pg_isready -h localhost -U %(ENV_SQL_USER)s -d %(ENV_SQL_DATABASE)s;
    do sleep 1;
  done;
  echo '[fastapi] PostgreSQL ready.';
  echo '[fastapi] Waiting for MongoDB...';
  until mongosh --quiet --eval 'db.runCommand({ping:1})' > /dev/null 2>&1;
    do sleep 1;
  done;
  echo '[fastapi] MongoDB ready. Starting uvicorn...';
  cd /home/user/app/backend && python3 -m uvicorn src.main:app --host 0.0.0.0 --port 7860
"
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SQL_HOST="localhost",MONGODB_URI="mongodb://localhost:27017/eventry",SQL_USER="admin",SQL_PASSWORD="password",SQL_DATABASE="eventry",PYTHONPATH="/home/user/app/backend",JWT_SECRET_KEY="xGLSZvOTb+KrlfjvNlovxg28tjp6Gj"