# supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=1
stdout_logfile=/tmp/postgres.log
stderr_logfile=/tmp/postgres.err

[program:mongodb]
command=mongod --dbpath /var/lib/mongodb --bind_ip_all
autostart=true
autorestart=true
priority=1
stdout_logfile=/tmp/mongodb.log
stderr_logfile=/tmp/mongodb.err

[program:fastapi]
command=python3 -m uvicorn src.main:app --host 0.0.0.0 --port 7860
directory=/home/user/app/backend
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/tmp/fastapi.log
stderr_logfile=/tmp/fastapi.err
environment=SQL_HOST="localhost",MONGODB_URI="mongodb://localhost:27017/eventry",SQL_USER="admin",SQL_PASSWORD="password",SQL_DATABASE="eventry"
